# systemd/.config/systemd/user/update_checker.service
[Unit]
Description=Check updates from AUR and Pacman

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/update_checker.sh
