# systemd/.config/systemd/user/update_checker.timer
[Unit]
Description=Run the updates checker every 30 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
